# scripts/com.bookmark-watcher.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.bookmark-watcher</string>

    <key>ProgramArguments</key>
    <array>
        <string>{{PROJECT_DIR}}/.venv/bin/python3</string>
        <string>{{PROJECT_DIR}}/watcher.py</string>
    </array>

    <key>WorkingDirectory</key>
    <string>{{PROJECT_DIR}}</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>{{PROJECT_DIR}}/logs/launchd.out.log</string>

    <key>StandardErrorPath</key>
    <string>{{PROJECT_DIR}}/logs/launchd.err.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>
</dict>
</plist>
